# deploy/systemd/virtual-org-control-api.service
# Systemd unit for keeping the AI Venture Studio control API running on a Linux server.
# Copy this file to /etc/systemd/system/virtual-org-control-api.service and adjust the paths
# if your repo lives somewhere other than /root/virtual-org.

[Unit]
Description=AI Venture Studio Control API
After=network.target postgresql.service
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/virtual-org
ExecStart=/root/virtual-org/.venv/bin/python3 /root/virtual-org/scripts/run_api.py
Restart=always
RestartSec=5

# Keep the process readable and predictable.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
